# deploy/alarm-web.service
[Unit]
Description=Emfit Alarm Web App
After=network.target emfit-qs2.service
Wants=emfit-qs2.service

[Service]
Type=simple
User=okazaki
WorkingDirectory=/home/okazaki/emfit-alarm/alarmapp
Environment=PATH=/home/okazaki/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/home/okazaki/.local/bin/uvicorn app:app --host 0.0.0.0 --port 8123
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target