# deploy/emfit-qs2.service
[Unit]
Description=Emfit QS2 API
After=network.target

[Service]
Type=simple
User=minipc
WorkingDirectory=/home/minipc/chromecast/emfit-api
EnvironmentFile=/home/minipc/chromecast/deploy/emfit.env
ExecStart=/usr/bin/python3 -m qs2.main
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
